# dot_config/autostart/768_tmux_client.desktop
[Desktop Entry]
Exec=/usr/bin/xterm -geometry 84x28 -class special-xterm -e "~/.tmux/plugins/rat_scripts/autostart.sh"
Name=tmux
Comment=A preconfigured tmux session for clients
Type=Application
