# dot_config/autostart/gammastep.desktop
[Desktop Entry]
Exec=/usr/bin/gammastep -r
Name=gammastep
Comment=set color temperature of display according to time of day
Type=Application
