# dot_config/autostart/greyscale.sh.desktop
[Desktop Entry]
Exec=/home/mango/.config/picom/greyscale.sh
Name=greyscale.sh
Type=Application
