# dot_config/autostart/keynav.desktop
[Desktop Entry]
Exec=/usr/bin/keynav "loadconfig ~/.config/keynav/.keynavrc"
Name=keynav
Type=Application
