# dot_config/autostart/numen.desktop
[Desktop Entry]
Exec=/usr/local/bin/numen
Name=numen
Type=Application
