# dot_config/autostart/snixembed.desktop
[Desktop Entry]
Exec=/usr/bin/snixembed
Name=snixembed
Comment=proxy the StatusNotifierItem protocol to the XEmbed System Tray
Type=Application
