# dot_config/autostart/talon.desktop
[Desktop Entry]
Exec=/home/mango/.talon/talon/run.sh
Name=run.sh
Type=Application
